# gunicorn.service
[Unit]
Description=gunicorn daemon
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/home/youruser/myproject
ExecStart=/usr/bin/gunicorn myproject.wsgi:application --bind 127.0.0.1:8000

[Install]
WantedBy=multi-user.target
